# crontab
# Fetch Tibber prices every day at hh:mm
* * * * * /usr/local/bin/python /app/fetch/fetch_tibber.py >> /app/data/fetch_tibber.log 2>&1

