# crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
TZ=Europe/Stockholm

# Tibber once per day (13:00 local)
00 13 * * * root cd /app && /usr/local/bin/python fetch/fetch_tibber.py >> /app/data/fetch_tibber.log 2>&1

# Calendar every 5 minutes
*/5 * * * * root cd /app && /usr/local/bin/python fetch/fetch_calendar.py >> /app/data/fetch_calendar.log 2>&1


